# deploy/systemd/home-automation.service
# Boot + stack bring-up. Install together with the API watchdog:
#   cd /opt/home-automation && sudo ./scripts/install-systemd.sh
#
# User=/Group=/HOME= are rewritten from SUDO_USER when using install-systemd.sh.
# Deploy user must be in the "docker" group.

[Unit]
Description=Home Automation stack (Docker Compose production)
Documentation=https://github.com/elijah286/kerry_home_automation
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/home-automation
User=elijah286
Group=elijah286
Environment=HOME=/home/elijah286
# .env must exist here for compose (secrets, POSTGRES_PASSWORD, etc.)
ExecStart=/usr/bin/docker compose -f /opt/home-automation/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /opt/home-automation/docker-compose.prod.yml stop
TimeoutStartSec=0
# If compose fails, systemd marks the unit failed (visible in journalctl)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
