# deploy/systemd/home-automation-watchdog.service
# Runs scripts/watchdog-stack.sh — self-heal when :3000/api/health fails.
# Install with the main stack: sudo ./scripts/install-systemd.sh
#
# Logs: journalctl -u home-automation-watchdog.service -f

[Unit]
Description=Home Automation API watchdog (docker compose up -d if unhealthy)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=elijah286
Group=elijah286
Environment=HOME=/home/elijah286
WorkingDirectory=/opt/home-automation
ExecStart=/opt/home-automation/scripts/watchdog-stack.sh

[Install]
WantedBy=multi-user.target
