# deploy/systemd/home-automation-watchdog.timer
# Runs the watchdog on a short interval and once shortly after boot.
# Install: sudo ./scripts/install-systemd.sh (copies this unit and enables the timer).

[Unit]
Description=Timer for Home Automation API watchdog

[Timer]
OnBootSec=90s
OnUnitActiveSec=2min
AccuracySec=30s
RandomizedDelaySec=20
Persistent=true

[Install]
WantedBy=timers.target
