# deploy/systemd/home-automation.service
# Boot + stack bring-up. Install together with the API watchdog:
#   cd /opt/home-automation && sudo ./scripts/install-systemd.sh
#
# User=/Group=/HOME= are rewritten from SUDO_USER when using install-systemd.sh.
# Deploy user must be in the "docker" group.

[Unit]
Description=Home Automation stack (Docker Compose production)
Documentation=https://github.com/elijah286/kerry_home_automation
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/home-automation
User=elijah286
Group=elijah286
Environment=HOME=/home/elijah286
# Host-level provisioning BEFORE compose: installs intel-media-va-driver-non-free
# if missing, writes correct RENDER_GID / VIDEO_GID into .env (so the go2rtc
# container's group_add matches this host's actual render group), and sanity-
# checks /dev/dri. The `+` prefix runs this as root even though the main
# ExecStart runs as the deploy user. Failures log to journalctl but are
# non-fatal — the stack still comes up even if apt-get hiccups.
ExecStartPre=+/opt/home-automation/scripts/host-prereqs.sh
# .env must exist here for compose (secrets, POSTGRES_PASSWORD, etc.)
ExecStart=/usr/bin/docker compose -f /opt/home-automation/docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f /opt/home-automation/docker-compose.prod.yml stop
TimeoutStartSec=0
# If compose fails, systemd marks the unit failed (visible in journalctl)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
